# conf/supervisor/redis.conf
[program:redis]
command=service redis start
user=root
autostart=true
autorestart=true
startsecs=3